# chef/cookbooks/fort/files/default/postgres/upstart.conf
description "PostgreSQL database server"

start on runlevel [2345]
stop on runlevel [!2345]

setuid postgres
setgid postgres

respawn

exec postgres -XD /var/pgsql/data
